# pipeline_carga_dados/crontab

* * * * * python /home/project/pipeline.py >> /var/log/cron.log 2>&1