# pipeline_carga_dados/crontab

0 12 * * * /usr/local/bin/python /home/project/pipeline.py >> /var/log/cron.log 2>&1